# docker/supervisord.combined.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:gateway]
command=/opt/venv/bin/python -m hermes_cli.main gateway run --replace
directory=/opt/hermes/hermes-agent
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=HERMES_HOME="%(ENV_HERMES_HOME)s",PATH="/opt/venv/bin:%(ENV_PATH)s"

[program:webui]
command=/opt/venv/bin/python /opt/hermes/hermes-webui/server.py
directory=/opt/hermes/hermes-webui
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=HERMES_HOME="%(ENV_HERMES_HOME)s",HERMES_BASE_HOME="%(ENV_HERMES_BASE_HOME)s",HERMES_CONFIG_PATH="%(ENV_HERMES_CONFIG_PATH)s",HERMES_WEBUI_STATE_DIR="%(ENV_HERMES_WEBUI_STATE_DIR)s",HERMES_WEBUI_DEFAULT_WORKSPACE="%(ENV_HERMES_WEBUI_DEFAULT_WORKSPACE)s",HERMES_WEBUI_AGENT_DIR="%(ENV_HERMES_WEBUI_AGENT_DIR)s",HERMES_WEBUI_PYTHON="%(ENV_HERMES_WEBUI_PYTHON)s",HERMES_WEBUI_HOST="%(ENV_HERMES_WEBUI_HOST)s",HERMES_WEBUI_PORT="%(ENV_HERMES_WEBUI_PORT)s",PATH="/opt/venv/bin:%(ENV_PATH)s"
