# deploy/karmabox-pro/systemd/hermes-gateway.service.example
[Unit]
Description=Hermes Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/karmabox_pro/env.systemd
WorkingDirectory=/opt/hermes/hermes-agent
ExecStart=/opt/hermes/hermes-agent/.venv/bin/python -m hermes_cli.main gateway run --replace
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
